# roles/install/files/gitlab-runner.ini
[program:gitlab-runner]
command=/home/$USER/bin/gitlab-runner.sh
